# prometheus/custom_exporter/bash_custom_exporter.service
[Unit]
Description=Custom User Metrics Exporter
Documentation=https://github.com/prometheus/node_exporter
After=network.target
Wants=network.target

[Service]
Type=simple
User=custom_exporter
Group=custom_exporter
ExecStart=/bin/bash /usr/local/bin/user_sessions_exporter.sh
Restart=on-failure
RestartSec=5s
TimeoutStopSec=10s
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=custom_exporter

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes

[Install]
WantedBy=multi-user.target
